# dotfiles/config/systemd/user/vdirsyncer.service.d/override.conf
[Unit]
ConditionPathExists=/dev/yubikey

[Service]
Environment=DISPLAY=:0.0
